# ebindkeys-zipit/files/z2-ebindkeys
#!/bin/sh /etc/rc.common
# Start Zipit Z2 ebindkeys Daemon

START=19

start() {
	/usr/sbin/ebindkeys -f /etc/ebindkeysrc
	/usr/sbin/ebindkeys -f /etc/ebindpwrbtn
}

stop() {
	killall ebindkeys
}
