# bldaemon-zipit/files/z2-bldaemon
#!/bin/sh /etc/rc.common
# Start Zipit Z2 Backlight Daemon

START=11

start() {
	/usr/sbin/bldaemon &
}

stop() {
	killall bldaemon
}
